# contrib/apt-cron/apt.cron
#!/bin/sh

[ ! -f /var/lock/subsys/apt ] && exit 0

[ -f /etc/sysconfig/apt ] && . /etc/sysconfig/apt

[ $CHECK_ONLY == "yes" ] && OPTS="$OPTS --check-only"
OPTS="$OPTS $EXTRA_OPTIONS"

export HOME=/root
/usr/bin/apt-get -q update && /usr/bin/apt-get -q -y $OPTS dist-upgrade
